# arch/system-maintenance/system-maintenance.timer
[Unit]
Description=Daily system maintenance

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target